# fishy-backend.service
[Unit]
Description=Fishy Backend Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/fishy
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

